# server/crontab
37 */3 * * *	python /app/scripts/cleanRuntime.py
